# archlinux/ctrl-app.service
[Unit]
Description=BBB Control Application

[Service]
ExecStart=/usr/bin/ctrl-app

[Install]
WantedBy=multi-user.target

