# archlinux/tc.service
[Unit]
Description=BBB Control Application Enhanced

[Service]
ExecStart=/usr/bin/tcdr

[Install]
WantedBy=multi-user.target

